# config/systemd/sidekiq@.service.erb
[Unit]
Description=Sidekiq worker for <%= fetch(:application) %> (<%= fetch(:stage) %> / <%= fetch(:rails_env) %>) (instance %i)
After=syslog.target network.target

[Service]
Type=notify
WatchdogSec=10

Environment = PATH=<%= fetch(:rbenv_path) %>/shims:/usr/local/bin:/usr/bin:/bin
Environment = PWD=<%= current_path %>
Environment = RAILS_ENV=<%= fetch(:rails_env) %>
Environment = SIDEKIQ_WORKER=%i
EnvironmentFile=<%= current_path %>/.env

WorkingDirectory=<%= current_path %>

ExecStart=/home/<%= fetch(:user) %>/.rbenv/bin/rbenv exec bundle exec sidekiq -e <%= fetch(:rails_env) %> -C <%= current_path %>/config/sidekiq/<%= fetch(:stage) %>/%i.yml -g <%= fetch(:application) %>-sidekiq-<%= fetch(:stage) %>-%i

# Use `systemctl kill -s TSTP ` to quiet the Sidekiq process
UMask=0002
# ExecStop=/bin/kill -TERM $MAINPID

RestartSec=5
Restart=on-failure

TimeoutSec=900
TimeoutStopSec=30

SyslogIdentifier=<%= fetch(:application) %>-sidekiq-worker-<%= fetch(:stage) %>-%i

StandardOutput=append:<%= File.join(shared_path, 'log', 'sidekiq.log') %>
StandardError=append:<%= File.join(shared_path, 'log', 'sidekiq.log') %>

[Install]
WantedBy=default.target
